# supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=Xvfb :1 -screen 0 1280x800x24
autorestart=true

[program:xfce]
command=startxfce4
environment=DISPLAY=":1"
autorestart=true

[program:vnc]
command=x11vnc -display :1 -forever -shared -rfbport 5901
autorestart=true

[program:novnc]
command=websockify 6080 localhost:5901 --web=/opt/novnc
autorestart=true

[program:ttyd]
command=ttyd -p 7681 bash
autorestart=true

[program:codeserver]
command=code-server --bind-addr 0.0.0.0:8081 --auth none /workspace
autorestart=true

[program:filebrowser]
command=filebrowser -r /workspace -p 8082
autorestart=true

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
