# supervisord.conf
; supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=info

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x800x24
autorestart=true
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err

[program:dbus]
command=/usr/bin/dbus-daemon --system --fork
autorestart=true
stdout_logfile=/var/log/dbus.log
stderr_logfile=/var/log/dbus.err

[program:xfce]
command=/bin/bash -lc "export DISPLAY=:1 && startxfce4"
environment=DISPLAY=":1"
autorestart=true
stdout_logfile=/var/log/xfce.log
stderr_logfile=/var/log/xfce.err

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -forever -shared -rfbport 5901 -rfbauth /home/dev/.vnc/passwd
autorestart=true
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err

[program:novnc]
command=/usr/local/bin/websockify --web=/opt/novnc 6080 localhost:5901
autorestart=true
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.err

[program:ttyd]
command=/usr/local/bin/ttyd -p 7681 -c dev:dev bash
autorestart=true
stdout_logfile=/var/log/ttyd.log
stderr_logfile=/var/log/ttyd.err

[program:codeserver]
command=/usr/bin/code-server --bind-addr 0.0.0.0:8080 --auth password --password ${CODE_PASSWORD} ${WORKDIR}
directory=/home/dev
user=dev
autorestart=true
stdout_logfile=/var/log/code-server.log
stderr_logfile=/var/log/code-server.err

[program:filebrowser]
command=/usr/local/bin/filebrowser -r ${WORKDIR} -p 8090
user=dev
autorestart=true
stdout_logfile=/var/log/filebrowser.log
stderr_logfile=/var/log/filebrowser.err
